# .config/autostart/pidgin.desktop
[Desktop Entry]
Type=Application
Exec=pidgin
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_US]=pidgin
Name=pidgin
Comment[en_US]=
Comment=
